# config/autostart/tutanota-desktop.desktop
[Desktop Entry]
	Type=Application
	Version=3.96.13
	Name=tutanota-desktop
	Comment=tutanota-desktop startup script
	Exec=/opt/tutanota-desktop/tutanota-desktop -a
	StartupNotify=false
	Terminal=false